# transcode.service
[Unit]
Description=Starts Transcode
[Service]
ExecStart=/usr/bin/transcode.sh
[Install]
WantedBy=multi-user.target
