# transcode.service
[Unit]
Description=Starts Transcode
[Service]
Type=oneshot
ExecStart=/usr/bin/transcode.sh
